# config/supervisor.d/cron.ini
[program:cron]
command = crond -f
autostart=true
autorestart=true
